# scripts/com.crunchtime.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.crunchtime.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mcotse/.nvm/versions/node/v25.3.0/bin/node</string>
    <string>/Users/mcotse/Developer/crunchtime/dist/server/index.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DB_PATH</key>
    <string>/Users/mcotse/Developer/crunchtime/crunchtime.db</string>
    <key>CF_TEAM_DOMAIN</key>
    <string>chummycollective</string>
    <key>PORT</key>
    <string>4000</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/mcotse/Developer/crunchtime</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/mcotse/Library/Logs/crunchtime-server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mcotse/Library/Logs/crunchtime-server.log</string>
</dict>
</plist>
